# signal-crawler/deploy/signal-crawler.service
[Unit]
Description=Signal crawler — scan HN, funding, and job boards for lead signals
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/deployer/north-cloud
ExecStart=/usr/bin/docker compose -f docker-compose.base.yml -f docker-compose.prod.yml run --rm signal-crawler
StandardOutput=journal
StandardError=journal
