# signal-crawler/deploy/signal-crawler.timer
[Unit]
Description=Signal crawler daily scan

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
